# com.leszek.analizator-growth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.leszek.analizator-growth</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_ROOT}}/venv/bin/python</string>
        <string>{{PROJECT_ROOT}}/app.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <dict>
        <key>NetworkState</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/analizator-growth.log</string>
    
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/analizator-growth-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>FLASK_ENV</key>
        <string>production</string>
        <key>PYTHONPATH</key>
        <string>{{PROJECT_ROOT}}</string>
    </dict>
</dict>
</plist> 